# roles/kubelet/files/docker-gc.service
[Service]
Type=oneshot
ExecStart=-/bin/bash -c "docker rmi $$(docker images -aq --filter dangling=true)"
