# roles/kubelet/files/docker-gc.timer
[Timer]
OnActiveSec=15min

[Install]
WantedBy=timers.target
